# roles/backup/templates/restic-restore-test.timer.j2
[Unit]
Description=Run restic monthly restore test

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target


