# roles/backup/templates/restic-backup.timer.j2
[Unit]
Description=Run restic backup daily

[Timer]
OnCalendar={{ backup_timer_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target


