# roles/backup/templates/restic-prune.timer.j2
[Unit]
Description=Run restic forget+prune daily

[Timer]
OnCalendar={{ prune_timer_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target


